# crontab
0	2	*	*	*	/home/ec2-user/user-repo/yslow.sh
*/3	*	*	*	*	/home/ec2-user/user-repo/yslow.sh new
